# src/gdalcubes/supervisord.conf
[supervisord]
nodaemon=true
logfile=/opt/supervisord.log

[program:gdalcubes_server]
command=gdalcubes_server -p 1111 -t 2
autorestart=unexpected
exitcodes=0
startsecs=0
